# etc/init.d/bluetooth-obexd-openrc
#!/sbin/openrc-run

description="Starts the Bluetooth OBEX daemon"
supervisor="supervise-daemon"
command="/usr/libexec/bluetooth/obexd"

depend() {
	need bluetooth
	provide obexd
}
