# etc/init.d/bluetooth-openrc
#!/sbin/openrc-run
# Copyright (c) 2007-2021 The Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

description="Starts the Bluetooth daemon"
command="/usr/sbin/bluetoothd"
pidfile="/run/bluetoothd.pid"

depend() {
	# Needs the D-Bus message bus service
	need dbus
	# Should start after kernel modules are loaded to ensure
	# Bluetooth hardware drivers are available.
	after modules
	provide bluetooth
}
