# etc/init.d/bluez-openrc
#!/sbin/openrc-run

description="Start Bluetooth daemons"
description_start="Starting Bluetooth daemons"
description_stop="Stopping Bluetooth daemons"

command="/usr/sbin/bluetoothd"
pidfile="/run/bluetoothd.pid"

depend() {
    need mountall
    need dbus
    after bootmisc
}

start() {
    ebegin "${description_start}"
    start-stop-daemon --start --exec "${command}" --pidfile "${pidfile}" --background
    eend $?
}

stop() {
    ebegin "${description_stop}"
    start-stop-daemon --stop --exec "${command}" --pidfile "${pidfile}"
    eend $?
}
