# etc/init.d/bluetooth-meshd-openrc
#!/sbin/openrc-run

description="Starts the Bluetooth Mesh daemon"
supervisor="supervise-daemon"
command="/usr/libexec/bluetooth/bluetooth-meshd"

depend() {
	need bluetooth
	provide meshd
}
