# etc/init.d/bluetooth-meshd-openrc
#!/sbin/openrc-run
# Copyright (c) 2007-2021 The Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

description="Starts the Bluetooth Mesh daemon"
command="/usr/libexec/bluetooth/bluetooth-meshd"
command_args="${meshd_opts}"
pidfile="/run/bluetooth-meshd.pid"

depend() {
	# Needs the D-Bus message bus service
	need dbus
	# Needs the main bluetooth service to be running first
	need bluetooth-openrc
}

start_pre() {
	# These are good checks to ensure the system is ready for the daemon.
	if [ ! -d /sys/class/bluetooth ]; then
		eerror "Bluetooth hardware support not found (/sys/class/bluetooth)."
		return 1
	fi

	if [ ! -x "${command}" ]; then
		eerror "${command} is not installed or not executable."
		return 1
	fi
}
